# streamlit_service.service
[Unit]
Description=Student Code Debugging Platform - Streamlit Service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/student-platform
Environment="PATH=/home/ubuntu/student-platform/venv/bin"
ExecStart=/home/ubuntu/student-platform/venv/bin/streamlit run app.py --server.port=8501 --server.address=0.0.0.0 --server.headless=true
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
